# jbidwatcher.conf
[program:jbidwatcher]
command='/opt/jbidwatcher/jbidwatcher.sh'
environment=JAVA_HOME="/usr/lib/jvm/java-8-oracle",DISPLAY=":0.0"
directory=/opt/jbidwatcher
autorestart=true
